# adum_scrapper.service
[Unit]
Description=ADUM Inscription Script Scrapper

[Service]
WorkingDirectory=
EnvironmentFile=
ExecStart=./entrypoint.sh
Type=oneshot
